# crontab
@daily      root    /opt/rsync-daily.sh
