# crontab
@daily      root    /opt/backup-daily.sh
@monthly    root    /opt/backup-monthly.sh
